# etc/cron.d/dcec.cron
SHELL=/bin/sh
MAILTO=m.sokolov@haulmont.com
0   7   *   *   *   zabbix  /etc/zabbix/scripts/dcec/dcec_update_cache.py >/dev/null 2>/dev/null
