# etc/cron.d/zdcec.cron
SHELL=/bin/sh
MAILTO=m.sokolov@haulmont.com
30   6   *   *   *   zabbix  /usr/bin/python3 /etc/zabbix/scripts/zdcec/zdcec_update_cache.py >/dev/null 2>/dev/null
